# roles/guacamole/templates/guacd.service.j2
[Unit]
Description=Guacamole Proxy Daemon
After=network.target

[Service]
ExecStart=/usr/local/sbin/guacd -f
Restart=on-failure

[Install]
WantedBy=multi-user.target
